# deploy/systemd/b2-verification.service
[Unit]
Description=B2 deterministic-failure short-circuit verification (one-shot)

[Service]
Type=oneshot
ExecStart=/home/patrick/dev/Shadow/scripts/b2_verification.sh
StandardOutput=journal
StandardError=journal
